# crontab
*/30    *       *       *       *       root    python /root/XX-XXXXXXXXXX.py
