# crontab
*/30    *       *       *       *       root    python /root/28-01186e2bbeff.py
